# crontab
* * * * * if ! systemctl is-active --quiet caddy; then systemctl restart caddy; fi
